# flask-restful-restx/supervisor-test.ini
[program:TaxPlatformService]
directory=/root/projects/TaxPlatformService
command=sh start.sh
autostart=true
autorestart=true
stderr_logfile=/root/projects/TaxPlatformService/logs/test_stderr.log
stdout_logfile=/root/projects/TaxPlatformService/logs/test_stdout.log
environment=LOG_PATH="/root/logs/TaxPlatformService",RETURNS_PATH="/root/data/TaxPlatformService/returns"