# flask-restful-restx/supervisor-production.ini
[program:TaxPlatformService]
directory=/root/projects/TaxPlatformService
command=sh start.sh
autostart=true
autorestart=true
stderr_logfile=/root/logs/TaxPlatformService/test_stderr.log
stdout_logfile=/root/logs/TaxPlatformService/test_stdout.log
environment=MYSQL_HOST="127.0.0.1",MYSQL_PASSWORD="Chuangke888#",REDIS_HOST="127.0.0.1",DEBUG="0",LOG_PATH="/root/logs/TaxPlatformService",FASTDFS_CONF_PATH="/usr/local/fastdfs/conf/client.conf"